# old/ub/packages/salt-minion/files/salt-call.cron
## This file is managed by Salt
RANDOM_DELAY=10

*/15 * * * * root /usr/bin/salt-call state.apply &>/dev/null

