# old/ub/packages/apt-check/files/apt-check.cron
## This file is managed by Salt"
#!/bin/bash

/opt/webhost/bin/apt-check  > /dev/null 2>&1
