# states/pkgs/salt-minion/files/salt-call.cron
### andy.dustin@gmail.com [rev: 9c0d006]
## This file is managed by Salt
RANDOM_DELAY=10

*/15 * * * * root nice salt-call --master=salt.finalduty.me state.apply &>/dev/null

